# Configuration/etc/profile.d/gantibg.sh
#! /bin/bash
./usr/bin/gd-bg --image /usr/share/anemona/bg.png